# alpine/xlx/s6_overlay_WORKING/etc/s6-overlay/s6-rc.d/httpd/run
#!/command/with-contenv sh

exec /usr/sbin/httpd -D FOREGROUND
