# ambed/docker/s6_overlay/etc/s6-overlay/s6-rc.d/ambed/run
#!/command/with-contenv sh

IP=$( hostname -I )

exec /ambed/ambed ${IP}
